# .ebextensions/crontab
* * * * * /opt/python/run/venv/bin/python3.6 /opt/python/current/app/ingest.py



